# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid kenpachi
chdir /home/kenpachi/sites/SITENAME/source

exec /home/kenpachi/.pyenv/versions/django/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    superlists.wsgi:application
